# docs/examples/supervisor/mediasnake.conf
[program:mediasnake]
command = uwsgi_python -H env --socket data/uwsgi.sock --mount=/mediasnake=mediasnake/wsgi.py
directory = /srv/mediasnake
user = www-data
